# config/supervisor.conf
[program:frame_page]
command=/usr/local/bin/docsify serve docs
directory=/var/www/frame_pages
restart=always
user=root
